# overlayfs/etc/rc.local
#!/bin/sh

## rc.local style script for running custom stuff during boot
## place commands before "exit 0"

exit 0

